# openwrt/base-files/etc/rc.local
#!/bin/sh 
#-*-sh-*-
#
# $Id: rc.local $
#
# Author: Markus Stenberg <markus stenberg@iki.fi>
#
# Copyright (c) 2012 cisco Systems, Inc.
#
# Created:       Wed Oct 31 11:32:33 2012 mstenber
# Last modified: Fri Oct 11 06:44:45 2013 mstenber
# Edit time:     44 min
#

SETUPONCE=/root/localsetup.sh

if [ -f $SETUPONCE ]
then
    . $SETUPONCE
    rm -f $SETUPONCE
    if [ ! -f $SETUPONCE ]
    then
        sync # paranoia?
        reboot
    fi
fi

# Start bird+pm.lua
#/etc/init.d/hnet start

